# .config/systemd/user/scrot.service
[Unit]
Description=Screenshot
Documentation=man:scrot(1)

[Service]
Type=oneshot
ExecStart=/usr/bin/scrot %h/.shots/%%FT%%T.png
